# references/deploy/app.service.example
# systemd — copiar para deploy/app.service
# Ajustar User, WorkingDirectory e domínio conforme project_config vps.*

[Unit]
Description=SaaS Next.js App
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/meu-saas
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/var/www/meu-saas/.env
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
